# workspace/cron/session-governance.cron
# 会话文件治理 - 每2小时执行
*/10 * * * * /root/.openclaw/workspace/scripts/session-cleanup-governor.sh

# Gateway内存监控 - 每5分钟执行
*/5 * * * * /root/.openclaw/workspace/scripts/gateway-memory-governor.sh

# EventBus 阈值扫描 - 每10分钟执行
*/10 * * * * cd /root/.openclaw/workspace && /usr/bin/node infrastructure/event-bus/sensors/threshold-scanner.js >> logs/threshold-scanner.log 2>&1
